# debian/postinst
#!/bin/bash

detect_desktop_environment() {
  if [ -n "$XDG_CURRENT_DESKTOP" ]; then
    echo "$XDG_CURRENT_DESKTOP"
  elif [ -n "$DESKTOP_SESSION" ]; then
    echo "$DESKTOP_SESSION"
  else
    if [ -d "/usr/share/xsessions" ]; then
      sessions=$(ls /usr/share/xsessions/ | sed 's/\.desktop$//')
      if [ -n "$sessions" ]; then
        echo "$sessions" | head -n 1
      else
        echo "cinnamon"
      fi
    else
      echo "cinnamon"
    fi
  fi
}

case "$1" in
    configure)

      # Get Real User ID
      RUID=$(who | awk 'FNR == 1 {print $1}')
      # Get Real User UID
      RUSER_UID=$(id -u ${RUID})

      if [ "$RUID" != "" ] ; then
        if [ ! -f /home/${RUID}/.config/autostart/hamonikr-theme.desktop ] ; then
          [ ! -d "/home/${RUID}/.config/autostart/" ] && mkdir -p /home/${RUID}/.config/autostart/
          cp /etc/skel/.config/autostart/hamonikr-theme.desktop /home/${RUID}/.config/autostart/
          chown ${RUID}:${RUID} /home/${RUID}/.config/autostart -R
        fi

        if [ ! -f /home/${RUID}/.config/autostart/plank.desktop ] ; then
          [ ! -d "/home/${RUID}/.config/autostart/" ] && mkdir -p /home/${RUID}/.config/autostart/
          cp /etc/skel/.config/autostart/plank.desktop /home/${RUID}/.config/autostart/
          chown ${RUID}:${RUID} /home/${RUID}/.config/autostart -R

          [ ! -d "/home/${RUID}/.config/plank/dock1/launchers/" ] && mkdir -p /home/${RUID}/.config/plank/dock1/launchers/
          cp /etc/skel/.config/plank/dock1/launchers/*.dockitem /home/${RUID}/.config/plank/dock1/launchers/
          chown ${RUID}:${RUID} /home/${RUID}/.config/plank -R          
        fi

        if [ ! -f /home/${RUID}/.config/autostart/default-conky.desktop ] ; then
          [ ! -d "/home/${RUID}/.config/autostart/" ] && mkdir -p /home/${RUID}/.config/autostart/
          cp /etc/skel/.config/autostart/default-conky.desktop /home/${RUID}/.config/autostart/
          chown ${RUID}:${RUID} /home/${RUID}/.config/autostart -R
        fi        
      fi

      if [ ! -z "$RUID" ] ; then
          sudo -u ${RUID} DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/${RUSER_UID}/bus" dconf load /org/gnome/terminal/legacy/profiles:/:$(gsettings list-recursively org.gnome.Terminal.ProfilesList | grep default | cut -d' ' -f3 | tr -d "'")/ < /usr/share/hamonikr-themes/gnome-terminal.dconf
      fi

      # For menu icon change
      if [ -d "/home/${RUID}/.config/cinnamon/spices/menu@cinnamon.org/" ] ; then
        rm -rf "/home/${RUID}/.config/cinnamon/spices/menu@cinnamon.org/*.json"
        cp /usr/share/cinnamon/applets/menu@cinnamon.org/menu.json /home/${RUID}/.config/cinnamon/spices/menu@cinnamon.org/0.json
      fi

      gtk-update-icon-cache || true
      fc-cache -f

      if which systemctl >/dev/null 2>&1 ; then
          systemctl enable hamonikr-themes.service
          systemctl daemon-reload
      fi      

      # https://forum.whale.naver.com/topic/46349/
      if [ "$(dpkg-query -W --showformat='${db:Status-Status}' naver-whale-stable 2>&1)" = "installed" ] ; then
        if [ ! -f /etc/apt/trusted.gpg.d/whale-key.gpg ] ; then
          apt-key export EF6C07F6 | sudo gpg --dearmour -o /etc/apt/trusted.gpg.d/whale-key.gpg
          echo "Fixed Naver whale GPG KEY warning."
        fi
      fi

      # Detect Desktop Environment
      DESKTOP_ENV=$(detect_desktop_environment)
      DESKTOP_ENV=$(echo "$DESKTOP_ENV" | tr '[:upper:]' '[:lower:]')

      # Update lightdm configuration
      if [ -f /etc/lightdm/lightdm.conf.d/99-hamonikr-greeter.conf ]; then
        sed -i "s/^user-session=.*/user-session=${DESKTOP_ENV}/" /etc/lightdm/lightdm.conf.d/99-hamonikr-greeter.conf
      fi   

    ;;
    abort-upgrade|abort-remove|abort-deconfigure)

    ;;
    triggered)
    ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

exit 0